# ananicy-rules/debian/postinst
#!/bin/sh

set -e

systemctl enable ananicy-cpp --now
systemctl start ananicy-cpp